# modules/lldp/etc/uci-defaults/uconfig-lldp

SHADOW=/etc/uconfig/shadow
touch $SHADOW/lldpd
