# modules/tailscale/etc/uci-defaults/uconfig-tailscale
#!/bin/sh

cp /etc/config/tailscale /etc/uconfig/shadow
